# var/lib/dpkg/info/courier-imap-ssl.list
/.
/usr
/usr/lib
/usr/lib/courier
/usr/lib/courier/mkimapdcert
/usr/share
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/mkimapdcert.8.gz
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/courier-imap-ssl
/usr/share/doc
/usr/share/doc/courier-imap-ssl
/usr/share/doc/courier-imap-ssl/copyright
/usr/sbin
/etc
/etc/init.d
/etc/init.d/courier-imap-ssl
/etc/logcheck
/etc/logcheck/ignore.d.server
/etc/logcheck/ignore.d.server/courier-imap-ssl
/etc/logcheck/violations.ignore.d
/etc/logcheck/violations.ignore.d/courier-imap-ssl
/etc/courier
/etc/courier/imapd.cnf
/etc/courier/imapd-ssl
/usr/lib/courier/imapd.pem
/usr/share/man/man8/imapd-ssl.8.gz
/usr/share/doc/courier-imap-ssl/NEWS.Debian.gz
/usr/share/doc/courier-imap-ssl/changelog.Debian.gz
/usr/sbin/imapd-ssl
/usr/sbin/mkimapdcert
