# var/lib/dpkg/info/landscape-client.conffiles
/etc/logrotate.d/landscape-client
/etc/init.d/landscape-client
