# var/lib/dpkg/info/supervisor.list
/.
/var
/var/log
/var/log/supervisor
/usr
/usr/bin
/usr/bin/supervisord
/usr/bin/pidproxy
/usr/bin/supervisorctl
/usr/bin/echo_supervisord_conf
/usr/lib
/usr/lib/python2.7
/usr/lib/python2.7/dist-packages
/usr/lib/python2.7/dist-packages/supervisor-3.0b2.egg-info
/usr/lib/python2.7/dist-packages/supervisor
/usr/lib/python2.7/dist-packages/supervisor/medusa
/usr/lib/python2.7/dist-packages/supervisor/medusa/thread
/usr/lib/python2.7/dist-packages/supervisor/medusa/test
/usr/lib/python2.7/dist-packages/supervisor/medusa/docs
/usr/lib/python2.7/dist-packages/supervisor/medusa/demo
/usr/lib/python2.7/dist-packages/supervisor/scripts
/usr/lib/python2.7/dist-packages/supervisor/tests
/usr/lib/python2.7/dist-packages/supervisor/tests/fixtures
/usr/lib/python2.7/dist-packages/supervisor/skel
/usr/lib/python2.7/dist-packages/supervisor/ui
/usr/lib/python2.7/dist-packages/supervisor/ui/images
/usr/lib/python2.7/dist-packages/supervisor/ui/stylesheets
/usr/share
/usr/share/pyshared
/usr/share/pyshared/supervisor-3.0b2.egg-info
/usr/share/pyshared/supervisor-3.0b2.egg-info/top_level.txt
/usr/share/pyshared/supervisor-3.0b2.egg-info/entry_points.txt
/usr/share/pyshared/supervisor-3.0b2.egg-info/not-zip-safe
/usr/share/pyshared/supervisor-3.0b2.egg-info/dependency_links.txt
/usr/share/pyshared/supervisor-3.0b2.egg-info/SOURCES.txt
/usr/share/pyshared/supervisor-3.0b2.egg-info/requires.txt
/usr/share/pyshared/supervisor-3.0b2.egg-info/PKG-INFO
/usr/share/pyshared/supervisor-3.0b2.egg-info/namespace_packages.txt
/usr/share/pyshared/supervisor
/usr/share/pyshared/supervisor/supervisorctl.py
/usr/share/pyshared/supervisor/version.txt
/usr/share/pyshared/supervisor/events.py
/usr/share/pyshared/supervisor/dispatchers.py
/usr/share/pyshared/supervisor/xmlrpc.py
/usr/share/pyshared/supervisor/loggers.py
/usr/share/pyshared/supervisor/options.py
/usr/share/pyshared/supervisor/http.py
/usr/share/pyshared/supervisor/medusa
/usr/share/pyshared/supervisor/medusa/monitor_client.py
/usr/share/pyshared/supervisor/medusa/resolver.py
/usr/share/pyshared/supervisor/medusa/counter.py
/usr/share/pyshared/supervisor/medusa/TODO.txt
/usr/share/pyshared/supervisor/medusa/medusa_gif.py
/usr/share/pyshared/supervisor/medusa/status_handler.py
/usr/share/pyshared/supervisor/medusa/default_handler.py
/usr/share/pyshared/supervisor/medusa/producers.py
/usr/share/pyshared/supervisor/medusa/redirecting_handler.py
/usr/share/pyshared/supervisor/medusa/setup.py
/usr/share/pyshared/supervisor/medusa/xmlrpc_handler.py
/usr/share/pyshared/supervisor/medusa/chat_server.py
/usr/share/pyshared/supervisor/medusa/thread
/usr/share/pyshared/supervisor/medusa/thread/thread_channel.py
/usr/share/pyshared/supervisor/medusa/thread/select_trigger.py
/usr/share/pyshared/supervisor/medusa/thread/test_module.py
/usr/share/pyshared/supervisor/medusa/thread/thread_handler.py
/usr/share/pyshared/supervisor/medusa/thread/pi_module.py
/usr/share/pyshared/supervisor/medusa/event_loop.py
/usr/share/pyshared/supervisor/medusa/auth_handler.py
/usr/share/pyshared/supervisor/medusa/Makefile
/usr/share/pyshared/supervisor/medusa/virtual_handler.py
/usr/share/pyshared/supervisor/medusa/monitor_client_win32.py
/usr/share/pyshared/supervisor/medusa/asynchat_25.py
/usr/share/pyshared/supervisor/medusa/filesys.py
/usr/share/pyshared/supervisor/medusa/test
/usr/share/pyshared/supervisor/medusa/test/tests.txt
/usr/share/pyshared/supervisor/medusa/test/test_medusa.py
/usr/share/pyshared/supervisor/medusa/test/max_sockets.py
/usr/share/pyshared/supervisor/medusa/test/test_producers.py
/usr/share/pyshared/supervisor/medusa/test/test_single_11.py
/usr/share/pyshared/supervisor/medusa/test/asyn_http_bench.py
/usr/share/pyshared/supervisor/medusa/test/bench.py
/usr/share/pyshared/supervisor/medusa/test/test_11.py
/usr/share/pyshared/supervisor/medusa/test/test_lb.py
/usr/share/pyshared/supervisor/medusa/ftp_server.py
/usr/share/pyshared/supervisor/medusa/script_handler.py
/usr/share/pyshared/supervisor/medusa/unix_user_handler.py
/usr/share/pyshared/supervisor/medusa/docs
/usr/share/pyshared/supervisor/medusa/docs/composing_producers.gif
/usr/share/pyshared/supervisor/medusa/docs/producers.gif
/usr/share/pyshared/supervisor/medusa/docs/programming.html
/usr/share/pyshared/supervisor/medusa/docs/async_blurbs.txt
/usr/share/pyshared/supervisor/medusa/docs/proxy_notes.txt
/usr/share/pyshared/supervisor/medusa/docs/data_flow.html
/usr/share/pyshared/supervisor/medusa/docs/tkinter.txt
/usr/share/pyshared/supervisor/medusa/docs/threads.txt
/usr/share/pyshared/supervisor/medusa/docs/debugging.txt
/usr/share/pyshared/supervisor/medusa/docs/data_flow.gif
/usr/share/pyshared/supervisor/medusa/docs/README.html
/usr/share/pyshared/supervisor/medusa/http_date.py
/usr/share/pyshared/supervisor/medusa/MANIFEST
/usr/share/pyshared/supervisor/medusa/rpc_client.py
/usr/share/pyshared/supervisor/medusa/INSTALL.txt
/usr/share/pyshared/supervisor/medusa/README.txt
/usr/share/pyshared/supervisor/medusa/put_handler.py
/usr/share/pyshared/supervisor/medusa/CHANGES.txt
/usr/share/pyshared/supervisor/medusa/asyncore_25.py
/usr/share/pyshared/supervisor/medusa/http_server.py
/usr/share/pyshared/supervisor/medusa/logger.py
/usr/share/pyshared/supervisor/medusa/__init__.py
/usr/share/pyshared/supervisor/medusa/m_syslog.py
/usr/share/pyshared/supervisor/medusa/rpc_server.py
/usr/share/pyshared/supervisor/medusa/monitor.py
/usr/share/pyshared/supervisor/medusa/demo
/usr/share/pyshared/supervisor/medusa/demo/winFTPserver.py
/usr/share/pyshared/supervisor/medusa/demo/simple_anon_ftpd.py
/usr/share/pyshared/supervisor/medusa/demo/start_medusa.py
/usr/share/pyshared/supervisor/medusa/demo/script_server.py
/usr/share/pyshared/supervisor/medusa/demo/publish.py
/usr/share/pyshared/supervisor/states.py
/usr/share/pyshared/supervisor/rpcinterface.py
/usr/share/pyshared/supervisor/scripts
/usr/share/pyshared/supervisor/scripts/sample_exiting_eventlistener.py
/usr/share/pyshared/supervisor/scripts/sample_eventlistener.py
/usr/share/pyshared/supervisor/scripts/sample_commevent.py
/usr/share/pyshared/supervisor/scripts/loop_listener.py
/usr/share/pyshared/supervisor/scripts/loop_eventgen.py
/usr/share/pyshared/supervisor/socket_manager.py
/usr/share/pyshared/supervisor/tests
/usr/share/pyshared/supervisor/tests/test_rpcinterfaces.py
/usr/share/pyshared/supervisor/tests/test_confecho.py
/usr/share/pyshared/supervisor/tests/fixtures
/usr/share/pyshared/supervisor/tests/fixtures/donothing.conf
/usr/share/pyshared/supervisor/tests/fixtures/spew.py
/usr/share/pyshared/supervisor/tests/fixtures/unkillable_spew.py
/usr/share/pyshared/supervisor/tests/fixtures/fakeos.py
/usr/share/pyshared/supervisor/tests/test_supervisorctl.py
/usr/share/pyshared/supervisor/tests/test_childutils.py
/usr/share/pyshared/supervisor/tests/test_dispatchers.py
/usr/share/pyshared/supervisor/tests/test_events.py
/usr/share/pyshared/supervisor/tests/test_loggers.py
/usr/share/pyshared/supervisor/tests/test_socket_manager.py
/usr/share/pyshared/supervisor/tests/test_options.py
/usr/share/pyshared/supervisor/tests/test_xmlrpc.py
/usr/share/pyshared/supervisor/tests/test_web.py
/usr/share/pyshared/supervisor/tests/test_http.py
/usr/share/pyshared/supervisor/tests/trackrefs.py
/usr/share/pyshared/supervisor/tests/test_states.py
/usr/share/pyshared/supervisor/tests/test_supervisord.py
/usr/share/pyshared/supervisor/tests/__init__.py
/usr/share/pyshared/supervisor/tests/test_datatypes.py
/usr/share/pyshared/supervisor/tests/base.py
/usr/share/pyshared/supervisor/tests/test_process.py
/usr/share/pyshared/supervisor/childutils.py
/usr/share/pyshared/supervisor/confecho.py
/usr/share/pyshared/supervisor/web.py
/usr/share/pyshared/supervisor/supervisord.py
/usr/share/pyshared/supervisor/http_client.py
/usr/share/pyshared/supervisor/datatypes.py
/usr/share/pyshared/supervisor/skel
/usr/share/pyshared/supervisor/skel/sample.conf
/usr/share/pyshared/supervisor/pidproxy.py
/usr/share/pyshared/supervisor/process.py
/usr/share/pyshared/supervisor/ui
/usr/share/pyshared/supervisor/ui/status.html
/usr/share/pyshared/supervisor/ui/tail.html
/usr/share/pyshared/supervisor/ui/images
/usr/share/pyshared/supervisor/ui/images/button_restart.gif
/usr/share/pyshared/supervisor/ui/images/supervisor.gif
/usr/share/pyshared/supervisor/ui/images/icon.png
/usr/share/pyshared/supervisor/ui/images/rule.gif
/usr/share/pyshared/supervisor/ui/images/state3.gif
/usr/share/pyshared/supervisor/ui/images/button_stop.gif
/usr/share/pyshared/supervisor/ui/images/state0.gif
/usr/share/pyshared/supervisor/ui/images/state2.gif
/usr/share/pyshared/supervisor/ui/images/button_refresh.gif
/usr/share/pyshared/supervisor/ui/images/state1.gif
/usr/share/pyshared/supervisor/ui/stylesheets
/usr/share/pyshared/supervisor/ui/stylesheets/supervisor.css
/usr/share/pyshared/supervisor-3.0b2-nspkg.pth
/usr/share/doc
/usr/share/doc/supervisor
/usr/share/doc/supervisor/TODO.txt
/usr/share/doc/supervisor/changelog.Debian.gz
/usr/share/doc/supervisor/README.Debian
/usr/share/doc/supervisor/copyright
/usr/share/doc/supervisor/README.rst
/usr/share/doc/supervisor/examples
/usr/share/doc/supervisor/examples/sample_exiting_eventlistener.py
/usr/share/doc/supervisor/examples/sample_eventlistener.py
/usr/share/doc/supervisor/examples/sample_commevent.py
/usr/share/doc/supervisor/examples/loop_listener.py
/usr/share/doc/supervisor/examples/sample.conf.gz
/usr/share/doc/supervisor/examples/loop_eventgen.py
/usr/share/python
/usr/share/python/ns
/usr/share/python/ns/supervisor
/etc
/etc/supervisor
/etc/supervisor/conf.d
/etc/supervisor/supervisord.conf
/etc/init.d
/etc/init.d/supervisor
/etc/default
/etc/default/supervisor
/usr/lib/python2.7/dist-packages/supervisor-3.0b2.egg-info/top_level.txt
/usr/lib/python2.7/dist-packages/supervisor-3.0b2.egg-info/entry_points.txt
/usr/lib/python2.7/dist-packages/supervisor-3.0b2.egg-info/not-zip-safe
/usr/lib/python2.7/dist-packages/supervisor-3.0b2.egg-info/dependency_links.txt
/usr/lib/python2.7/dist-packages/supervisor-3.0b2.egg-info/SOURCES.txt
/usr/lib/python2.7/dist-packages/supervisor-3.0b2.egg-info/requires.txt
/usr/lib/python2.7/dist-packages/supervisor-3.0b2.egg-info/PKG-INFO
/usr/lib/python2.7/dist-packages/supervisor-3.0b2.egg-info/namespace_packages.txt
/usr/lib/python2.7/dist-packages/supervisor/supervisorctl.py
/usr/lib/python2.7/dist-packages/supervisor/version.txt
/usr/lib/python2.7/dist-packages/supervisor/events.py
/usr/lib/python2.7/dist-packages/supervisor/dispatchers.py
/usr/lib/python2.7/dist-packages/supervisor/xmlrpc.py
/usr/lib/python2.7/dist-packages/supervisor/loggers.py
/usr/lib/python2.7/dist-packages/supervisor/options.py
/usr/lib/python2.7/dist-packages/supervisor/http.py
/usr/lib/python2.7/dist-packages/supervisor/medusa/monitor_client.py
/usr/lib/python2.7/dist-packages/supervisor/medusa/resolver.py
/usr/lib/python2.7/dist-packages/supervisor/medusa/counter.py
/usr/lib/python2.7/dist-packages/supervisor/medusa/TODO.txt
/usr/lib/python2.7/dist-packages/supervisor/medusa/medusa_gif.py
/usr/lib/python2.7/dist-packages/supervisor/medusa/status_handler.py
/usr/lib/python2.7/dist-packages/supervisor/medusa/default_handler.py
/usr/lib/python2.7/dist-packages/supervisor/medusa/producers.py
/usr/lib/python2.7/dist-packages/supervisor/medusa/redirecting_handler.py
/usr/lib/python2.7/dist-packages/supervisor/medusa/setup.py
/usr/lib/python2.7/dist-packages/supervisor/medusa/xmlrpc_handler.py
/usr/lib/python2.7/dist-packages/supervisor/medusa/chat_server.py
/usr/lib/python2.7/dist-packages/supervisor/medusa/thread/thread_channel.py
/usr/lib/python2.7/dist-packages/supervisor/medusa/thread/select_trigger.py
/usr/lib/python2.7/dist-packages/supervisor/medusa/thread/test_module.py
/usr/lib/python2.7/dist-packages/supervisor/medusa/thread/thread_handler.py
/usr/lib/python2.7/dist-packages/supervisor/medusa/thread/pi_module.py
/usr/lib/python2.7/dist-packages/supervisor/medusa/event_loop.py
/usr/lib/python2.7/dist-packages/supervisor/medusa/auth_handler.py
/usr/lib/python2.7/dist-packages/supervisor/medusa/Makefile
/usr/lib/python2.7/dist-packages/supervisor/medusa/virtual_handler.py
/usr/lib/python2.7/dist-packages/supervisor/medusa/monitor_client_win32.py
/usr/lib/python2.7/dist-packages/supervisor/medusa/asynchat_25.py
/usr/lib/python2.7/dist-packages/supervisor/medusa/filesys.py
/usr/lib/python2.7/dist-packages/supervisor/medusa/test/tests.txt
/usr/lib/python2.7/dist-packages/supervisor/medusa/test/test_medusa.py
/usr/lib/python2.7/dist-packages/supervisor/medusa/test/max_sockets.py
/usr/lib/python2.7/dist-packages/supervisor/medusa/test/test_producers.py
/usr/lib/python2.7/dist-packages/supervisor/medusa/test/test_single_11.py
/usr/lib/python2.7/dist-packages/supervisor/medusa/test/asyn_http_bench.py
/usr/lib/python2.7/dist-packages/supervisor/medusa/test/bench.py
/usr/lib/python2.7/dist-packages/supervisor/medusa/test/test_11.py
/usr/lib/python2.7/dist-packages/supervisor/medusa/test/test_lb.py
/usr/lib/python2.7/dist-packages/supervisor/medusa/ftp_server.py
/usr/lib/python2.7/dist-packages/supervisor/medusa/script_handler.py
/usr/lib/python2.7/dist-packages/supervisor/medusa/unix_user_handler.py
/usr/lib/python2.7/dist-packages/supervisor/medusa/docs/composing_producers.gif
/usr/lib/python2.7/dist-packages/supervisor/medusa/docs/producers.gif
/usr/lib/python2.7/dist-packages/supervisor/medusa/docs/programming.html
/usr/lib/python2.7/dist-packages/supervisor/medusa/docs/async_blurbs.txt
/usr/lib/python2.7/dist-packages/supervisor/medusa/docs/proxy_notes.txt
/usr/lib/python2.7/dist-packages/supervisor/medusa/docs/data_flow.html
/usr/lib/python2.7/dist-packages/supervisor/medusa/docs/tkinter.txt
/usr/lib/python2.7/dist-packages/supervisor/medusa/docs/threads.txt
/usr/lib/python2.7/dist-packages/supervisor/medusa/docs/debugging.txt
/usr/lib/python2.7/dist-packages/supervisor/medusa/docs/data_flow.gif
/usr/lib/python2.7/dist-packages/supervisor/medusa/docs/README.html
/usr/lib/python2.7/dist-packages/supervisor/medusa/http_date.py
/usr/lib/python2.7/dist-packages/supervisor/medusa/MANIFEST
/usr/lib/python2.7/dist-packages/supervisor/medusa/rpc_client.py
/usr/lib/python2.7/dist-packages/supervisor/medusa/INSTALL.txt
/usr/lib/python2.7/dist-packages/supervisor/medusa/README.txt
/usr/lib/python2.7/dist-packages/supervisor/medusa/put_handler.py
/usr/lib/python2.7/dist-packages/supervisor/medusa/CHANGES.txt
/usr/lib/python2.7/dist-packages/supervisor/medusa/asyncore_25.py
/usr/lib/python2.7/dist-packages/supervisor/medusa/http_server.py
/usr/lib/python2.7/dist-packages/supervisor/medusa/logger.py
/usr/lib/python2.7/dist-packages/supervisor/medusa/__init__.py
/usr/lib/python2.7/dist-packages/supervisor/medusa/m_syslog.py
/usr/lib/python2.7/dist-packages/supervisor/medusa/rpc_server.py
/usr/lib/python2.7/dist-packages/supervisor/medusa/monitor.py
/usr/lib/python2.7/dist-packages/supervisor/medusa/demo/winFTPserver.py
/usr/lib/python2.7/dist-packages/supervisor/medusa/demo/simple_anon_ftpd.py
/usr/lib/python2.7/dist-packages/supervisor/medusa/demo/start_medusa.py
/usr/lib/python2.7/dist-packages/supervisor/medusa/demo/script_server.py
/usr/lib/python2.7/dist-packages/supervisor/medusa/demo/publish.py
/usr/lib/python2.7/dist-packages/supervisor/states.py
/usr/lib/python2.7/dist-packages/supervisor/rpcinterface.py
/usr/lib/python2.7/dist-packages/supervisor/scripts/sample_exiting_eventlistener.py
/usr/lib/python2.7/dist-packages/supervisor/scripts/sample_eventlistener.py
/usr/lib/python2.7/dist-packages/supervisor/scripts/sample_commevent.py
/usr/lib/python2.7/dist-packages/supervisor/scripts/loop_listener.py
/usr/lib/python2.7/dist-packages/supervisor/scripts/loop_eventgen.py
/usr/lib/python2.7/dist-packages/supervisor/socket_manager.py
/usr/lib/python2.7/dist-packages/supervisor/tests/test_rpcinterfaces.py
/usr/lib/python2.7/dist-packages/supervisor/tests/test_confecho.py
/usr/lib/python2.7/dist-packages/supervisor/tests/fixtures/donothing.conf
/usr/lib/python2.7/dist-packages/supervisor/tests/fixtures/spew.py
/usr/lib/python2.7/dist-packages/supervisor/tests/fixtures/unkillable_spew.py
/usr/lib/python2.7/dist-packages/supervisor/tests/fixtures/fakeos.py
/usr/lib/python2.7/dist-packages/supervisor/tests/test_supervisorctl.py
/usr/lib/python2.7/dist-packages/supervisor/tests/test_childutils.py
/usr/lib/python2.7/dist-packages/supervisor/tests/test_dispatchers.py
/usr/lib/python2.7/dist-packages/supervisor/tests/test_events.py
/usr/lib/python2.7/dist-packages/supervisor/tests/test_loggers.py
/usr/lib/python2.7/dist-packages/supervisor/tests/test_socket_manager.py
/usr/lib/python2.7/dist-packages/supervisor/tests/test_options.py
/usr/lib/python2.7/dist-packages/supervisor/tests/test_xmlrpc.py
/usr/lib/python2.7/dist-packages/supervisor/tests/test_web.py
/usr/lib/python2.7/dist-packages/supervisor/tests/test_http.py
/usr/lib/python2.7/dist-packages/supervisor/tests/trackrefs.py
/usr/lib/python2.7/dist-packages/supervisor/tests/test_states.py
/usr/lib/python2.7/dist-packages/supervisor/tests/test_supervisord.py
/usr/lib/python2.7/dist-packages/supervisor/tests/__init__.py
/usr/lib/python2.7/dist-packages/supervisor/tests/test_datatypes.py
/usr/lib/python2.7/dist-packages/supervisor/tests/base.py
/usr/lib/python2.7/dist-packages/supervisor/tests/test_process.py
/usr/lib/python2.7/dist-packages/supervisor/childutils.py
/usr/lib/python2.7/dist-packages/supervisor/confecho.py
/usr/lib/python2.7/dist-packages/supervisor/web.py
/usr/lib/python2.7/dist-packages/supervisor/supervisord.py
/usr/lib/python2.7/dist-packages/supervisor/http_client.py
/usr/lib/python2.7/dist-packages/supervisor/datatypes.py
/usr/lib/python2.7/dist-packages/supervisor/skel/sample.conf
/usr/lib/python2.7/dist-packages/supervisor/pidproxy.py
/usr/lib/python2.7/dist-packages/supervisor/process.py
/usr/lib/python2.7/dist-packages/supervisor/ui/status.html
/usr/lib/python2.7/dist-packages/supervisor/ui/tail.html
/usr/lib/python2.7/dist-packages/supervisor/ui/images/button_restart.gif
/usr/lib/python2.7/dist-packages/supervisor/ui/images/supervisor.gif
/usr/lib/python2.7/dist-packages/supervisor/ui/images/icon.png
/usr/lib/python2.7/dist-packages/supervisor/ui/images/rule.gif
/usr/lib/python2.7/dist-packages/supervisor/ui/images/state3.gif
/usr/lib/python2.7/dist-packages/supervisor/ui/images/button_stop.gif
/usr/lib/python2.7/dist-packages/supervisor/ui/images/state0.gif
/usr/lib/python2.7/dist-packages/supervisor/ui/images/state2.gif
/usr/lib/python2.7/dist-packages/supervisor/ui/images/button_refresh.gif
/usr/lib/python2.7/dist-packages/supervisor/ui/images/state1.gif
/usr/lib/python2.7/dist-packages/supervisor/ui/stylesheets/supervisor.css
/usr/lib/python2.7/dist-packages/supervisor-3.0b2-nspkg.pth
