# var/lib/dpkg/info/vsftpd.conffiles
/etc/systemd/system/vsftpd.service
/etc/vsftpd.conf
/etc/logrotate.d/vsftpd
/etc/init/vsftpd.conf
/etc/pam.d/vsftpd
/etc/ftpusers
